# resources/rekoit-restore.service
[Unit]
Description=Restore REKOIT runtime after firmware update
After=home.mount
Wants=home.mount
Before=xochitl.service
ConditionPathExists=/home/root/rekoit/restore.sh

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh /home/root/rekoit/restore.sh

[Install]
WantedBy=multi-user.target
